# eduvulcan_token_fetcher/rootfs/etc/services.d/eduvulcan/finish
#!/usr/bin/env bashio
# ==============================================================================
# Stop container after the one-shot job finishes
# ==============================================================================

declare APP_EXIT_CODE=${1}

echo "${APP_EXIT_CODE}" > /run/s6-linux-init-container-results/exitcode
if [[ "${APP_EXIT_CODE}" -eq 0 ]]; then
  bashio::log.info "Token fetch completed. Stopping add-on."
else
  bashio::log.error "Token fetch failed with exit code ${APP_EXIT_CODE}. Stopping add-on."
fi

exec /run/s6/basedir/bin/halt
