# eduvulcan_token_fetcher/rootfs/etc/services.d/eduvulcan/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start EduVulcan token fetcher service
# ==============================================================================

exec /usr/bin/run.sh
