# example/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# ==============================================================================

## Add your code here

bashio::log.info "Hello world"

# exec my_program --with-params
